# OzricAddon/rootfs/etc/services.d/ozric/run
#!/usr/bin/with-contenv bashio

# Declare variables
declare message

## Get the 'message' key from the user config options.
message=$(bashio::config 'message')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${message:="Hello World..."}"

## Run your program
# exec my_program --with-params